# files/scripts/zramswap
### BEGIN INIT INFO
# Provides: zramswap
# Required-Start:
# Required-Stop:
# Default-Start: 2
# Default-Stop:
# Short-Description: Create swap partitions
# Description:
### END INIT INFO

# chmod +x /etc/init.d/zramswap
# update-rc.d actled defaults 2

. /lib/lsb/init-functions

case "$1" in
  start)
    log_daemon_msg "Creating swap partitions"
    zramswap start 2>/dev/null
    log_end_msg $?
    ;;
  *)
    echo "Usage: $0 start" >&2
    exit 3
    ;;
esac
